# hetrixtools-agent/rootfs/etc/s6-overlay/s6-rc.d/hetrixtools/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -eo pipefail

LIB=/usr/lib/hetrixtools
RUN_DIR=/run/hetrixtools

# --- validate SID (fail fast, clear message) ---
SID="$(bashio::config 'sid')"
if [[ ! "$SID" =~ ^[A-Za-z0-9]{32}$ ]]; then
	bashio::exit.nok "Invalid or missing 'sid' (need a 32-char HetrixTools Server ID)."
fi

# --- render cfg next to a working copy of the agent ---
mkdir -p "$RUN_DIR"
cp "$LIB/hetrixtools_agent.sh" "$RUN_DIR/hetrixtools_agent.sh"

export SID
export COLLECT_EVERY_SECONDS; COLLECT_EVERY_SECONDS="$(bashio::config 'collect_every_seconds')"
export CHECK_DRIVE_HEALTH;    CHECK_DRIVE_HEALTH="$(bashio::config 'check_drive_health' | grep -q true && echo 1 || echo 0)"
export CHECK_SOFT_RAID;       CHECK_SOFT_RAID="$(bashio::config 'check_soft_raid' | grep -q true && echo 1 || echo 0)"
export CHECK_REBOOT;          CHECK_REBOOT="$(bashio::config 'check_reboot' | grep -q true && echo 1 || echo 0)"
export RUNNING_PROCESSES;     RUNNING_PROCESSES="$(bashio::config 'running_processes' | grep -q true && echo 1 || echo 0)"
export CONNECTION_PORTS;      CONNECTION_PORTS="$(bashio::config 'connection_ports' | tr '\n' ',' | sed 's/,$//')"
export CHECK_SERVICES;        CHECK_SERVICES="$(bashio::config 'check_services' | tr '\n' ',' | sed 's/,$//')"

bash "$LIB/render-config.sh" "$RUN_DIR/hetrixtools.cfg"

if bashio::config.true 'dry_run'; then
	export DryRun=1
	bashio::log.warning "dry_run enabled: payloads printed, not sent."
fi

bashio::log.info "Starting HetrixTools agent loop (SID ${SID:0:6}…)."

# --- agent loop ---
# The agent self-paces each invocation to ~60s (it samples for a minute, then
# POSTs once). collect_every_seconds changes sample granularity, NOT the POST
# interval. The sleep below is only a crash-loop backstop: if a cycle exits
# early, backfill the rest of the minute so a broken agent can't hammer the API.
while true; do
	start=$(date +%s)
	bash "$RUN_DIR/hetrixtools_agent.sh" || bashio::log.warning "agent cycle exited non-zero"
	elapsed=$(( $(date +%s) - start ))
	# Guard against a backward clock step (NTP) yielding a negative duration.
	if [ "$elapsed" -lt 0 ]; then
		elapsed=0
	fi
	if [ "$elapsed" -lt 55 ]; then
		sleep $(( 60 - elapsed ))
	fi
done
